# roles/btc-rpc-explorer/templates/btc-rpc-explorer.service.j2



[Unit]
Description=btc-rpc-explorer

[Service]
Type=simple
ExecStart=btc-rpc-explorer --port 8080 
Restart=on-failure

[Install]
WantedBy=multi-user.target
